# test-service/files/test-service.init
#!/bin/sh /etc/rc.common
. /usr/lib/amx/scripts/amx_procd_init_functions.sh
USE_PROCD=1
START=96
STOP=10
name="test-service"
PROG="/usr/sbin/test-service"

start_service() {
    register_service ${name} ${PROG}
}

stop_service() {
    deregister_service ${name}
}
